# bot.conf
[program:telegrambot]
directory=/home/ubuntu/userbot
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu",USER="ubuntu"

stderr_logfile=/home/ubuntu/userbot/logfile_err.log
stdout_logfile=/home/ubuntu/userbot/logfile.log
